# pkg_files/deb/DEBIAN/postrm
update-rc.d maxgaugeair-agent remove
rm -f /usr/bin/maxgaugeair-agent
rm -rf /etc/maxgaugeair
rm -f /etc/init/maxgaugeair-agent.conf
rm -f /etc/init.d/maxgaugeair-agent

echo ""
echo "Agent for maxgauge Air uninstalled"
