# pkg_files/deb/DEBIAN/postrm
update-rc.d maxguageair remove
echo ""
echo "Agent for MaxGuage Air uninstalled"
